# .config/systemd/user/easyeffects.service
[Service]
ExecStart=/usr/bin/easyeffects --gapplication-service
Restart=always
RestartSec=5

